# btbms-display.service
[Unit]
Description=BtBmsDisplay Web Application
After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/j5_console/BtBmsDisplay
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStartPre=/usr/bin/npm install
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btbms-display

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/j5_console/BtBmsDisplay

[Install]
WantedBy=multi-user.target
