# btbms-display.service
[Unit]
Description=BtBmsDisplay Web Application
After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=seanfuchs
Group=seanfuchs
WorkingDirectory=/home/seanfuchs/BtBmsDisplay
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btbms-display

[Install]
WantedBy=multi-user.target
