# kiosk.service
[Unit]
Description=Chromium Kiosk Mode Service
After=network.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/Users/brettkettler/REPO/j5_console/BtBmsDisplay/start_kiosk.sh
Restart=always
RestartSec=5
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=HOME=/home/pi

[Install]
WantedBy=graphical-session.target