# kiosk.service
[Unit]
Description=Chromium Kiosk Mode Service
After=network.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/chmod +x /home/seanfuchs/Desktop/j5_console/BtBmsDisplay/start_kiosk.sh
ExecStart=/home/seanfuchs/Desktop/j5_console/BtBmsDisplay/start_kiosk.sh
Restart=always
RestartSec=5
User=seanfuchs
Group=seanfuchs
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/seanfuchs/.Xauthority
Environment=HOME=/home/seanfuchs

[Install]
WantedBy=graphical-session.target